# debian/rusty-socks/DEBIAN/postrm
#!/bin/sh
set -e

case "$1" in
    purge)
        # Stop and disable the service
        systemctl stop rusty-socks.service || true
        systemctl disable rusty-socks.service || true
        systemctl daemon-reload
        
        # Remove user and group
        if getent passwd rusty-socks >/dev/null; then
            deluser rusty-socks
        fi
        
        if getent group rusty-socks >/dev/null; then
            delgroup rusty-socks
        fi
        
        # Remove configuration and log directories
        rm -rf /etc/rusty-socks
        rm -rf /var/log/rusty-socks
        rm -rf /var/run/rusty-socks
        ;;
    
    remove)
        systemctl stop rusty-socks.service || true
        systemctl disable rusty-socks.service || true
        ;;
esac

# Automatically added by dh_installsystemd/13.11.4
if [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	systemctl --system daemon-reload >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.11.4
if [ "$1" = "purge" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper purge 'rusty-socks.service' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
