# debian/postrm
#!/bin/sh
set -e

case "$1" in
    purge)
        # Stop and disable the service
        systemctl stop rust-socksd.service || true
        systemctl disable rust-socksd.service || true
        systemctl daemon-reload

        # Remove user and group
        if getent passwd rust-socksd >/dev/null; then
            deluser rust-socksd
        fi

        if getent group rust-socksd >/dev/null; then
            delgroup rust-socksd
        fi

        # Remove configuration and log directories
        rm -rf /etc/rust-socksd
        rm -rf /var/log/rust-socksd
        rm -rf /var/lib/rust-socksd
        ;;

    remove)
        systemctl stop rust-socksd.service || true
        systemctl disable rust-socksd.service || true
        ;;
esac

#DEBHELPER#

exit 0
