# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Reload the systemd configuration
        systemd-sysusers
        systemctl daemon-reload

        # Set ownership and permissions
        chown -R rusty-socks:rusty-socks /var/lib/rusty-socks
        chown -R rusty-socks:rusty-socks /var/log/rusty-socks
        chown root:rusty-socks /etc/rusty-socks/config.yml
        chmod 640 /etc/rusty-socks/config.yml
        chown root:rusty-socks /etc/rusty-socks/users.yml
        chmod 660 /etc/rusty-socks/users.yml

        # Enable and start the service
        systemctl enable rusty-socks.service

        echo "Rusty SOCKS has been installed successfully."
        echo "Edit /etc/rusty-socks/config.yml to configure the server."
        echo "Start the service with: systemctl start rusty-socks"
        ;;
esac

#DEBHELPER#

exit 0
